# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[program:samba]
command=/samba-wrapper.sh
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/var/log/supervisor/samba-stdout.log
stderr_logfile=/var/log/supervisor/samba-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=100

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
